# deploy/okx-hft-grid.service
[Unit]
Description=OKX Altcoin HFT Grid Trading System
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=3
OnFailure=okx-hft-grid-failure@%n.service

[Service]
Type=notify
User=okxtrader
Group=okxtrader
WorkingDirectory=/opt/okx-hft-grid
ExecStart=/opt/okx-hft-grid/okx-hft-grid -config /etc/okx-hft-grid/config.yaml
EnvironmentFile=/etc/okx-hft-grid/env
StateDirectory=okx-hft-grid
StateDirectoryMode=0750
ProtectSystem=strict
ReadWritePaths=/var/lib/okx-hft-grid
ProtectHome=yes
NoNewPrivileges=yes
UMask=0077
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=okx-hft-grid
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
